# scripts/services/linux/myattention-watchdog.service
[Unit]
Description=IKE Runtime Watchdog
After=network.target ike-api.service
Wants=ike-api.service

[Service]
Type=simple
WorkingDirectory=/opt/ike
Environment=PYTHONUTF8=1
ExecStart=/opt/ike/.venv/bin/python /opt/ike/runtime_watchdog.py --mode distributed
Restart=always
RestartSec=5
StandardOutput=append:/var/log/ike/watchdog.log
StandardError=append:/var/log/ike/watchdog.log

[Install]
WantedBy=multi-user.target
